# systemd/media-actions-api.service
[Unit]
Description=Media Actions API (playlist add, track metadata, grouping)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SPAINIFY_USER__
WorkingDirectory=__SPAINIFY_ROOT__/apps/media-actions-api
EnvironmentFile=__SPAINIFY_ROOT__/apps/media-actions-api/.env
ExecStart=/usr/bin/node __SPAINIFY_ROOT__/apps/media-actions-api/server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
