# systemd/spotify_display.service
[Unit]
Description=Spotify Display Control
After=network.target weather-dashboard.service
Requires=weather-dashboard.service

[Service]
Type=simple
User=aspain
# We keep the venv under backend, but the code under apps/spotify-display
WorkingDirectory=/home/aspain/spainify
EnvironmentFile=-/home/aspain/spainify/apps/spotify-display/.env
ExecStartPre=/bin/sleep 30
ExecStart=/home/aspain/spainify/backend/venv/bin/python /home/aspain/spainify/apps/spotify-display/spotify_display_control.py
Environment="DISPLAY=:0"
Restart=always

[Install]
WantedBy=multi-user.target
