# systemd/spotify_display.service
[Unit]
Description=Spotify Display Control
After=network-online.target weather-dashboard.service sonify-serve.service sonos-http-api.service
Wants=network-online.target weather-dashboard.service sonify-serve.service sonos-http-api.service

[Service]
Type=simple
User=__SPAINIFY_USER__
# We keep the venv under backend, but the code under apps/spotify-display
WorkingDirectory=__SPAINIFY_ROOT__
EnvironmentFile=-__SPAINIFY_ROOT__/apps/spotify-display/.env
EnvironmentFile=-__SPAINIFY_ROOT__/apps/weather-dashboard/.env
EnvironmentFile=-__SPAINIFY_ROOT__/.spainify-device.env
ExecStartPre=__SPAINIFY_ROOT__/backend/venv/bin/python __SPAINIFY_ROOT__/apps/spotify-display/wait_for_ready_state.py
ExecStart=__SPAINIFY_ROOT__/backend/venv/bin/python __SPAINIFY_ROOT__/apps/spotify-display/spotify_display_control.py
Environment="DISPLAY=:0"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
